# Vault/vault.service
[Unit]
Description=Vault server

[Service]
Type=simple
ExecStart=/opt/VaultServer/Vault/vault server -config=/opt/VaultServer/Vault/config.hcl
